# onhandserver.service
[Unit]
Description=Telegram server resource monitor bot
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/opt/onhandserver
EnvironmentFile=/opt/onhandserver/.env
ExecStart=/opt/onhandserver/venv/bin/python /opt/onhandserver/bot.py
Restart=on-failure
RestartSec=10

[Install]
WantedBy=multi-user.target
